# airootfs/etc/skel/.config/autostart/xfce4-power-manager.desktop
[Desktop Entry]
Type=Application
Name=XFCE Power Manager
Exec=xfce4-power-manager
Icon=xfce4-power-manager
Terminal=false
Categories=System;
X-GNOME-Autostart-enabled=true
