# airootfs/etc/systemd/system/sshd.service.d/hardening.conf
# Systemd Service Hardening for SSH
# /etc/systemd/system/sshd.service.d/hardening.conf

[Service]
# Security hardening options

# Filesystem protections
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/empty

# Kernel protections
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Restrict namespaces
RestrictNamespaces=true

# Restrict realtime
RestrictRealtime=true

# Private /tmp
PrivateTmp=true

# Restrict address families
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Capability restrictions
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_CHOWN CAP_DAC_READ_SEARCH CAP_AUDIT_WRITE

# Lock personality
LockPersonality=true

# Remove IPC
RemoveIPC=true

# Protect hostname
ProtectHostname=true
