# airootfs/etc/systemd/system/setup-hunter.service
[Unit]
Description=Hunter OS First Boot Setup
After=sddm.service graphical.target
Wants=sddm.service
ConditionPathExists=!/var/lib/hunter-setup-done

[Service]
Type=oneshot
# Always mark as done even on failure to prevent infinite re-runs on next boot
ExecStart=/usr/local/bin/hunter-first-boot.sh
ExecStopPost=/usr/bin/touch /var/lib/hunter-setup-done
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
