# airootfs/etc/systemd/system/NetworkManager.service.d/hardening.conf
# Systemd Service Hardening for NetworkManager
# /etc/systemd/system/NetworkManager.service.d/hardening.conf

[Service]
# Security hardening options

# Filesystem protections
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/NetworkManager /run/NetworkManager /etc/NetworkManager

# Kernel protections
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Restrict namespaces
RestrictNamespaces=true

# Restrict realtime
RestrictRealtime=true

# No new privileges
NoNewPrivileges=false

# Private /tmp
PrivateTmp=true

# Restrict address families
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Capability restrictions
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID CAP_KILL

# Lock personality
LockPersonality=true

# Remove IPC
RemoveIPC=true
