# airootfs/etc/systemd/system/NetworkManager.service.d/hardening.conf
# Systemd Service Hardening for NetworkManager
# /etc/systemd/system/NetworkManager.service.d/hardening.conf

[Service]
# Balanced security hardening that doesn't break NetworkManager

# Filesystem protections (relaxed for NetworkManager to work)
ProtectSystem=false
ProtectHome=false

# Kernel protections
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=true
ProtectControlGroups=false

# Private /tmp
PrivateTmp=true

# Restrict address families (NetworkManager needs these)
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6 AF_PACKET

# Lock personality
LockPersonality=true
